# scrapyd-logparser/supervisord.d/logparser.ini
[program:logparser]
directory=/app
command = logparser -dir /app -ss 127.0.0.1:6800 --sleep 30 --verbose
autostart = true
startsecs = 5
autorestart = true
startretries = 3
redirect_stderr = true
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 5
stdout_logfile = /tmp/logparser.log